# crontab
1 * * * * * python /script.py >> /logfile
